# check-battery-user.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
