# check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/fred/.config/systemd/user/check_battery.sh
